# deploy/eridu-parser.service
[Unit]
Description=Eridu Inventory Parser
After=network.target

[Service]
Type=simple
User=eridu
Group=eridu
WorkingDirectory=/opt/eridu-ops-api/services/inventory_parser
Environment=PATH=/opt/eridu-ops-api/services/inventory_parser/.venv/bin
Environment=OMP_NUM_THREADS=4
Environment=MKL_NUM_THREADS=4
ExecStart=/opt/eridu-ops-api/services/inventory_parser/.venv/bin/gunicorn \
    --bind 127.0.0.1:5001 \
    --workers 1 \
    --preload \
    --timeout 300 \
    app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
