# deploy/eridu-bond100-sweep.service
[Unit]
Description=Eridu Bond 100 daily incremental /rank tail-fetch
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=eridu
Group=eridu
WorkingDirectory=/opt/eridu-ops-api/services/bond100
Environment=PATH=/opt/eridu-ops-api/services/bond100/.venv/bin
Environment=PYTHONUNBUFFERED=1
# Same DB the API serves from (outside the code dir so git pull never touches it).
Environment=BOND100_DB_PATH=/opt/eridu-ops-api/var/bond100.sqlite
# ARONA_TOKEN ('<email>:<token>') lives in this file, VPS only, so the token
# binds to the VPS IP. Leading "-" tolerates the file being absent.
EnvironmentFile=-/opt/eridu-ops-api/.env
# Incremental: reads the current bond-100 count and fetches only the newly-added
# records at the tail (~2-4 calls). Needs a prior `sweep_rank.py --global` to seed
# the store; re-run --global manually to resync after a removal or once arona fixes
# a broken record.
ExecStart=/opt/eridu-ops-api/services/bond100/.venv/bin/python sweep_rank.py --tail
